# workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/wwwroot/Sync/artisan schedule:run >> /dev/null 2>&1

